# supervisord.ini
[program:achtung5]
command = node server.js
directory = /home/janerist/apps/achtung5
environment = NODE_ENV="production",PORT_WWW="8080"
