# supervisord.ini
[program:achtung5]
command = node server.js
directory = /usr/local/apps/achtung5
environment = NODE_ENV=production
